# Library/LaunchDaemons/etc.hosts.relaxed.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>RunAtLoad</key>
    <true/>
    <key>Label</key>
    <string>etc.hosts.relaxed</string>
    <key>ProgramArguments</key>
    <array>
      <string>bash</string>
      <string>-c</string>
      <string>/bin/cp /etc/hosts.relaxed /etc/hosts; chmod a+r /etc/hosts</string>
    </array>
    <key>StandardOutPath</key>
    <string>/var/log/etc.hosts.relaxed</string>
    <key>StandardErrorPath</key>
    <string>/var/log/etc.hosts.relaxed</string>
    <key>StartCalendarInterval</key>
    <dict>
      <key>Hour</key>
      <integer>18</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
  </dict>
</plist>
